# allinone/tree-postinstall/etc/services/syslog-ng/run
#!/bin/sh
exec syslog-ng --no-caps -F -f /run/syslog-ng.conf \
    --persist-file=${DATA_DIR}/syslog-ng.persist \
    --pidfile /run/syslog-ng.pid
